# packages/os/linux/variants/elizaos-linux/config/includes.chroot/etc/systemd/system/elizaos-agent.service
[Unit]
Description=elizaOS local agent runtime (system-wide)
After=network-online.target elizaos-first-boot.service
Wants=network-online.target
Requires=elizaos-first-boot.service

[Service]
Type=simple
User=elizaos
Group=elizaos
WorkingDirectory=/var/lib/elizaos
Environment=ELIZA_STATE_DIR=/var/lib/elizaos
Environment=ELIZA_API_PORT=31337
ExecStart=/opt/elizaos/bin/bun /opt/elizaos/app/server.js --headless --port=31337
ExecStartPost=/usr/local/lib/elizaos/wait-agent-health.sh http://127.0.0.1:31337/api/health
Restart=on-failure
RestartSec=2
SyslogIdentifier=elizaos-agent
StandardOutput=journal+console
StandardError=journal+console
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/elizaos /var/log/elizaos
ProtectHome=yes

[Install]
WantedBy=multi-user.target
