# packages/os/linux/variants/elizaos-linux/config/includes.chroot/etc/systemd/system/elizaos-first-boot.service
[Unit]
Description=elizaOS first-boot bootstrap
After=local-fs.target
Before=elizaos-agent.service
ConditionPathExists=!/var/lib/elizaos/.first-boot-complete

[Service]
Type=oneshot
ExecStart=/usr/lib/elizaos/first-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
