# packages/os/linux/variants/milady-tails/tails/config/chroot_local-includes/etc/systemd/system/milady.service
[Unit]
Description=elizaOS AI desktop home supervisor
After=display-manager.service
ConditionPathExists=/usr/local/bin/milady
ConditionPathExists=/opt/milady/bin/launcher
RefuseManualStop=yes
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/lib/elizaos/milady-keeper
Restart=always
RestartSec=2
SyslogIdentifier=milady-supervisor
