# packages/os/linux/variants/milady-tails/tails/config/chroot_local-includes/etc/systemd/system/milady.path
[Unit]
Description=Start elizaOS when the live user session is ready

[Path]
PathExists=/run/user/1000/bus
Unit=milady.service

[Install]
WantedBy=multi-user.target
