# packages/os/linux/variants/milady-tails/tails/config/chroot_local-includes/etc/systemd/user/milady.service
[Unit]
Description=elizaOS AI desktop home
ConditionUser=1000
After=elizaos-renderer.service graphical-session.target tails-dump-user-env.service
Wants=elizaos-renderer.service tails-dump-user-env.service
RefuseManualStop=yes
StartLimitIntervalSec=0

[Service]
Type=simple
Environment=MILADY_SYSTEMD_SUPERVISED=1
Environment=MILADY_WAIT_FOR_LOCK=1
ExecStart=/usr/local/lib/elizaos/start-elizaos-browser-user
Restart=always
RestartSec=2
SyslogIdentifier=milady

[Install]
WantedBy=default.target
