# packages/os/linux/elizaos/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/autologin.conf
# elizaOS Linux — console autologin fallback for headless/no-GDM boots
# (e.g. serial-console debugging on riscv64). agetty --autologin skips PAM
# authentication, so it lands a shell as the live user without the
# "Authentication failure" that live-config's password-based getty autologin
# produced. Username matches GDM AutomaticLogin (etc/gdm3/daemon.conf).
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin user --noclear %I $TERM
