# packages/os/linux/variants/elizaos-debian-riscv64/config/includes.chroot/etc/systemd/system/elizaos-agent.service
# elizaOS RV64 variant userland-bootstrap (Wave 2B / first-boot scaffold)
[Unit]
Description=elizaOS local agent
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/elizaos
Environment=ELIZA_STATE_DIR=/var/lib/elizaos
ExecStart=/opt/elizaos/bin/elizaos start --headless --port=31337
Restart=on-failure
RestartSec=5
User=elizaos
Group=elizaos
SyslogIdentifier=elizaos-agent

[Install]
WantedBy=multi-user.target
