# packages/os/linux/variants/elizaos-debian-riscv64/config/includes.chroot/etc/systemd/system/elizaos-kokoro-audio-smoke.service
[Unit]
Description=elizaOS Kokoro llama.cpp audio smoke
After=elizaos-first-boot.service sound.target
Wants=sound.target

[Service]
Type=oneshot
ExecStart=/usr/lib/elizaos/kokoro-audio-smoke.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
